# .config/autostart/bitwarden.desktop
[Desktop Entry]
Type=Application
Name=Bitwarden
Exec=/opt/bitwarden/Bitwarden.AppImage
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
